# ops/crontab
MAILTO=<mail@to>
0 23 * * * chronic ~/backup-db.sh